# trunk/benchmark_XML/scripts/regression.cron.template
# $Id: regression.cron.template,v 1.2 2001-09-23 18:37:20 paklein Exp $
#
# NOTE: replace [path] before submitting
#
0 1 * * * /bin/csh -f [path to scripts]regression.csh > [path to script]regression.LOG

